# dropbear/etc/sv/dropbear/run
#!/bin/sh

exec /usr/sbin/dropbear -F -B
